# supervisor-bot.conf
[program:trading-bot]
directory=/root
command=python3 run_bot.py
user=root
autostart=true
autorestart=false
stderr_logfile=/root/bot-err.log
stdout_logfile=/root/bot-log.log